# templates/storage/docker-storage-setup.service.epp
<%- | $exec_start,
      $environment_file
| -%>
[Unit]
Description=Docker Storage Setup
After=cloud-init.service
Before=docker.service

[Service]
Type=oneshot
ExecStart=<%= $exec_start -%>
EnvironmentFile=-<%= $environment_file -%>

[Install]
WantedBy=multi-user.target
